# startup-docker-loki.service
[Unit]
After=network.target

[Service]
ExecStart=docker plugin enable loki

[Install]
WantedBy=default.target